# air1_all_zones_cv_time_features_package/systemd/user/air1-all-zones-person-counter-cam1.service
[Unit]
Description=AIR-1 All-Zones Person Counter cam1
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/air1_all_zones_cv_time_features_package
EnvironmentFile=-%h/air1_all_zones_cv_time_features_package/web_app/.env
Environment=PERSON_COUNT_CAMERA_ID=cam1
Environment=PYTHON_BIN=python3
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/air1_all_zones_cv_time_features_package/.python-packages:%h/air1_all_zones_cv_time_features_package
ExecStart=/usr/bin/env bash -lc 'exec bash run_person_counter.sh >> logs/person_counter_cam1.log 2>&1'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
