# air1_all_zones_cv_time_features_package/systemd/user/air1-all-zones-mqtt-aggregator.service
[Unit]
Description=AIR-1 All-Zones CV MQTT Aggregator
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/smart-i-lab-testbed/air1_all_zones_cv_time_features_package
EnvironmentFile=-%h/smart-i-lab-testbed/air1_all_zones_cv_time_features_package/web_app/.env
Environment=PYTHON_BIN=python3
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/smart-i-lab-testbed/air1_all_zones_cv_time_features_package/.python-packages:%h/smart-i-lab-testbed/air1_all_zones_cv_time_features_package
ExecStart=/usr/bin/env bash -lc 'exec "$PYTHON_BIN" -m air1_all_zones.occupancy_mqtt_aggregator >> logs/mqtt_aggregator.log 2>&1'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
