# zone5_cv_time_features_package/systemd/user/zone5-live-app.service
[Unit]
Description=Zone 5 Live App
After=network-online.target zone5-live-collector.service
Wants=zone5-live-collector.service

[Service]
Type=simple
WorkingDirectory=%h/zone5_cv_time_features_package
EnvironmentFile=-%h/zone5_cv_time_features_package/web_app/.env
Environment=HOST=127.0.0.1
Environment=PORT=8000
Environment=ZONE5_DATA_SOURCE=live
Environment=ZONE5_PRODUCTION_POINTER=model/production_run.txt
Environment=PYTHON_BIN=python3
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/zone5_cv_time_features_package/.python-packages:%h/zone5_cv_time_features_package
ExecStart=/usr/bin/env bash -lc 'exec bash run_live_app.sh >> logs/live_app.log 2>&1'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
