# zone5_cv_time_features_package/systemd/user/zone5-sen55-collector.service
[Unit]
Description=Zone 5 SEN55 Collector
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/smart-i-lab-testbed/zone5_cv_time_features_package
EnvironmentFile=-%h/smart-i-lab-testbed/zone5_cv_time_features_package/web_app/.env
Environment=PYTHON_BIN=python3
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/smart-i-lab-testbed/zone5_cv_time_features_package/.python-packages:%h/smart-i-lab-testbed/zone5_cv_time_features_package
ExecStart=/usr/bin/env bash -lc 'exec bash run_sen55_collector.sh >> logs/sen55_collector.log 2>&1'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
